# projects/cloudcurio-monorepo-pro-dash/apps/litellm-cloudcurio/systemd/litellm-stack.service
[Unit]
Description=Litellm
[Service]
Type=oneshot
